# target/linux/ipq806x/base-files/etc/init.d/cpufreq
#!/bin/sh /etc/rc.common

START=15
USE_PROCD=1

start_service() {
	local max_freq
	max_freq=$(awk '{print $NF}' /sys/devices/system/cpu/cpu0/cpufreq/scaling_available_frequencies)

	config_load cpufreq
	config_get governor general governor "ondemand"
	config_get minfreq general minfreq "600000"
	config_get maxfreq general maxfreq "$max_freq"
	config_get sampling_down_factor general sampling_down_factor "10"
	config_get up_threshold general up_threshold "50"

	for f in /sys/devices/system/cpu/cpufreq/policy*; do
		grep -q $governor "${f}/scaling_available_governors" && \
			echo -n "$governor" > "${f}/scaling_governor"

		grep -q $minfreq "${f}/scaling_available_frequencies" && \
			echo -n "$minfreq" > "${f}/scaling_min_freq"

		grep -q $maxfreq "${f}/scaling_available_frequencies" && \
			echo -n "$maxfreq" > "${f}/scaling_max_freq"

		[ "$governor" = "ondemand" ] && {
			echo -n "$sampling_down_factor" > "${f}/../ondemand/sampling_down_factor"
			echo -n "$up_threshold" > "${f}/../ondemand/up_threshold"
		}
	done
}

service_triggers() {
	procd_add_reload_trigger cpufreq
}
